# systemd/email-sender.timer
[Unit]
Description=Weather report service

[Timer]
OnBootSec=15min
OnCalendar=Mon..Fri 20:29
Persistent=true

[Install]
WantedBy=timers.target
